# vagrant/supervisor/serpisori.conf
[program:serpisori]
command = sudo -u www-data /vagrant/bin/start_tornado.sh
user = root
stdout_logfile = /var/log/supervisor/serpisori.log
redirect_stderr = true