# vagrant/supervisor/serpisori_debug.conf
[program:serpisori_debug]
command = sudo -u www-data /vagrant/bin/start_tornado_debug.sh
user = root
stdout_logfile = /var/log/supervisor/serpisori_debug.log
redirect_stderr = true